# Laravel/alpine/crontab
*/1 * * * * php /app/artisan schedule:run
*/1 * * * * sh /init.sh